# deploy_automation/crontab
0 4 * * MON /root/scripts/update_script.sh
@reboot /root/scripts/update_script.sh
0 0 * * * /root/scripts/cron_mon.sh